# var/lib/dpkg/info/dnsmasq.list
/.
/usr
/usr/share
/usr/share/doc
/etc
/etc/default
/etc/default/dnsmasq
/etc/dnsmasq.d
/etc/dnsmasq.d/README
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/dnsmasq.conf
/etc/init.d
/etc/init.d/dnsmasq
/etc/dnsmasq.conf
/etc/resolvconf
/etc/resolvconf/update.d
/etc/resolvconf/update.d/dnsmasq
/usr/share/doc/dnsmasq
