# var/lib/dpkg/info/dropbear.list
/.
/etc
/etc/dropbear
/etc/dropbear/log
/etc/dropbear/log/run
/etc/dropbear/run
/etc/init.d
/etc/init.d/dropbear
/usr
/usr/sbin
/usr/sbin/dropbear
/usr/lib
/usr/lib/dropbear
/usr/lib/dropbear/dropbearconvert
/usr/bin
/usr/bin/dropbearkey
/usr/bin/scp
/usr/bin/dbclient
/usr/dropbearconvert
/usr/share
/usr/share/initramfs-tools
/usr/share/initramfs-tools/conf-hooks.d
/usr/share/initramfs-tools/conf-hooks.d/dropbear
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/dropbear
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-premount
/usr/share/initramfs-tools/scripts/init-premount/devpts
/usr/share/initramfs-tools/scripts/init-premount/dropbear
/usr/share/initramfs-tools/scripts/init-bottom
/usr/share/initramfs-tools/scripts/init-bottom/dropbear
/etc/dropbear/log/main
