# var/lib/dpkg/info/lighttpd.list
/.
/usr
/usr/local
/usr/local/sbin
/usr/local/sbin/lighttpd-angel
/usr/local/sbin/lighttpd
/usr/local/lib
/usr/local/lib/mod_usertrack.la
/usr/local/lib/mod_fastcgi.la
/usr/local/lib/mod_rewrite.so
/usr/local/lib/mod_redirect.la
/usr/local/lib/mod_flv_streaming.so
/usr/local/lib/mod_expire.so
/usr/local/lib/mod_flv_streaming.la
/usr/local/lib/mod_proxy.so
/usr/local/lib/mod_compress.la
/usr/local/lib/mod_status.la
/usr/local/lib/mod_evhost.la
/usr/local/lib/mod_dirlisting.so
/usr/local/lib/mod_access.so
/usr/local/lib/mod_cgi.so
/usr/local/lib/mod_evhost.so
/usr/local/lib/mod_scgi.so
/usr/local/lib/mod_scgi.la
/usr/local/lib/mod_setenv.la
/usr/local/lib/mod_userdir.la
/usr/local/lib/mod_proxy.la
/usr/local/lib/mod_usertrack.so
/usr/local/lib/mod_cml.so
/usr/local/lib/mod_accesslog.so
/usr/local/lib/mod_extforward.so
/usr/local/lib/mod_redirect.so
/usr/local/lib/mod_secdownload.la
/usr/local/lib/mod_cgi.la
/usr/local/lib/mod_setenv.so
/usr/local/lib/mod_simple_vhost.so
/usr/local/lib/mod_staticfile.la
/usr/local/lib/mod_mysql_vhost.so
/usr/local/lib/mod_evasive.la
/usr/local/lib/mod_extforward.la
/usr/local/lib/mod_status.so
/usr/local/lib/mod_fastcgi.so
/usr/local/lib/mod_expire.la
/usr/local/lib/mod_trigger_b4_dl.so
/usr/local/lib/mod_evasive.so
/usr/local/lib/mod_magnet.so
/usr/local/lib/mod_ssi.so
/usr/local/lib/mod_cml.la
/usr/local/lib/mod_webdav.la
/usr/local/lib/mod_compress.so
/usr/local/lib/mod_accesslog.la
/usr/local/lib/mod_dirlisting.la
/usr/local/lib/mod_webdav.so
/usr/local/lib/mod_simple_vhost.la
/usr/local/lib/mod_staticfile.so
/usr/local/lib/mod_auth.so
/usr/local/lib/mod_alias.so
/usr/local/lib/mod_rrdtool.la
/usr/local/lib/mod_secdownload.so
/usr/local/lib/mod_access.la
/usr/local/lib/mod_indexfile.la
/usr/local/lib/mod_auth.la
/usr/local/lib/mod_ssi.la
/usr/local/lib/mod_indexfile.so
/usr/local/lib/mod_mysql_vhost.la
/usr/local/lib/mod_rewrite.la
/usr/local/lib/mod_trigger_b4_dl.la
/usr/local/lib/mod_alias.la
/usr/local/lib/mod_magnet.la
/usr/local/lib/mod_userdir.so
/usr/local/lib/mod_rrdtool.so
/usr/local/share
/usr/local/share/man
/usr/local/share/man/man8
/usr/local/share/man/man8/lighttpd.8.gz
/var
/var/log
/var/log/lighttpd
/var/log/lighttpd/error.log
/var/run
/var/run/lighttpd
/var/run/lighttpd/php-fastcgi.socket
/var/cache
/var/cache/lighttpd
/var/cache/lighttpd/uploads
/etc
/etc/init.d
/etc/init.d/lighttpd
/etc/lighttpd
/etc/lighttpd/lighttpd.conf
/etc/lighttpd/conf-available
/etc/lighttpd/conf-available/10-accesslog.conf
/etc/lighttpd/conf-available/10-status.conf
/etc/lighttpd/conf-available/10-expire.conf
/etc/lighttpd/conf-available/10-rrdtool.conf
/etc/lighttpd/conf-available/10-usertrack.conf
/etc/lighttpd/conf-available/10-simple-vhost.conf
/etc/lighttpd/conf-available/05-auth.conf
/etc/lighttpd/conf-available/10-ssl.conf
/etc/lighttpd/conf-available/90-debian-doc.conf
/etc/lighttpd/conf-available/10-cgi.conf
/etc/lighttpd/conf-available/15-fastcgi-php.conf
/etc/lighttpd/conf-available/10-no-www.conf
/etc/lighttpd/conf-available/10-evhost.conf
/etc/lighttpd/conf-available/10-dir-listing.conf
/etc/lighttpd/conf-available/README
/etc/lighttpd/conf-available/10-ssi.conf
/etc/lighttpd/conf-available/10-userdir.conf
/etc/lighttpd/conf-available/10-proxy.conf
/etc/lighttpd/conf-available/10-fastcgi.conf
/etc/lighttpd/conf-available/10-evasive.conf
/etc/lighttpd/conf-available/10-flv-streaming.conf
/etc/lighttpd/conf-enabled
