# var/lib/dpkg/info/syslog-ng.list
/.
/etc
/etc/logrotate.d
/etc/logrotate.d/syslog-ng
/etc/syslog-ng
/etc/syslog-ng/syslog-ng.conf
/etc/logcheck
/etc/logcheck/violations.ignore.d
/etc/logcheck/violations.ignore.d/syslog-ng
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/syslog-ng
/etc/logcheck/ignore.d.paranoid
/etc/logcheck/ignore.d.paranoid/syslog-ng
/etc/init.d
/etc/init.d/syslog-ng
/etc/default
/etc/default/syslog-ng
/usr
/usr/bin
/usr/bin/loggen
/usr/bin/update-patterndb
/usr/bin/pdbtool
/usr/share
/usr/share/man
/usr/share/man/man8/syslog-ng.8.gz
/usr/share/man/man5/syslog-ng.conf.5.gz
/usr/share/man/man1/pdbtool.1.gz
/usr/share/doc
/usr/share/doc/syslog-ng
/usr/share/doc/syslog-ng/changelog.Debian.gz
/usr/share/doc/syslog-ng/copyright
/usr/share/doc/syslog-ng/README.Debian
/usr/share/doc/syslog-ng/NEWS.gz
/usr/share/doc/syslog-ng/changelog.gz
/usr/sbin
/usr/sbin/syslog-ng
/usr/sbin/syslog-ng-ctl
